# rootfs/etc/s6-overlay/s6-rc.d/whisparr/finish
#!/usr/bin/execlineb -P

# stop execution if fails 3 times on 60 seconds interval
s6-permafailon 60 3 0-100 exit
