# rootfs/etc/s6-overlay/s6-rc.d/whisparr/run
#!/usr/bin/with-contenv sh

CONFIG_FILE=/config/config.xml

# create initial config file
if [[ ! -f "$CONFIG_FILE" ]]; then
  TMP_FILE="/tmp/$0.$$.$RANDOM"
  cat <<EOF >"$TMP_FILE"
<Config>
  <AnalyticsEnabled>False</AnalyticsEnabled>
</Config>
EOF

  install -m644 -o $PUID -g $PGID "$TMP_FILE" "$CONFIG_FILE"
elif [[ $(grep -c "<AnalyticsEnabled>" "$CONFIG_FILE") == 0 ]]; then
  # if analytics doenst exist in config, insert one
  sed -i "/<Config>/a\ \ <AnalyticsEnabled>False<\/AnalyticsEnabled>" "$CONFIG_FILE"
else
  # disable analytics
  sed -i "/<AnalyticsEnabled>/s/>.*<\//>False<\//" "$CONFIG_FILE"
fi

# change chdir
cd /config

# delete pid before start
rm -f whisparr.pid

# run with limited permissions
exec s6-setuidgid $PUID:$PGID /app/bin/Whisparr -nobrowser -data=/config
